# upstart/uim.conf
# /etc/init/uim.conf
description "Run UIM (for BT on pandaboard)"
author      "Eero af Heurlin <rambo@iki.fi>"

start on filesystem
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    exec /usr/local/sbin/uim
end script